# etc/systemd/system/bluetooth-agent.service
[Unit]
Description=ZeroAuth Bluetooth Agent

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python /home/alarm/media-pi/scripts/simple-agent.py
